# configs/testdaemon/testdaemon.service
[Unit]
Description=Test Daemon Service

[Service]
ExecStart=bash -c "date >> /home/matias/testdaemon.service.log"