# configs/testdaemon/testdaemon.timer
[Unit]
Description=Test Daemon Timer

[Timer]
Persistent=true
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target